# env/etc/circusd.conf
start on filesystem and net-device-up IFACE=lo
exec /usr/local/ap/venv/circus/bin/circusd --daemon /etc/circus.ini
